# launchd/com.cursor.feishu-bridge.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cursor.feishu-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.cursor/cursor-lark-bridge/daemon/cursor-lark-bridge-daemon</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>__HOME__/.cursor/cursor-lark-bridge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/.cursor/cursor-lark-bridge/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.cursor/cursor-lark-bridge/logs/launchd-stderr.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
